# llama-connector/llama-connector.service
[Unit]
Description=Llama Connector WebSocket Server
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/llama-connector
ExecStart=/home/ubuntu/llama-connector/venv/bin/python3 main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
